# laravel-scheduler.conf
[program:laravel-scheduler]
process_name=%(program_name)s_%(process_num)02d
command=php /c/laragon/www/HomeFinanceManager-new/artisan schedule:work
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/c/laragon/www/HomeFinanceManager-new/storage/logs/scheduler.log
stopwaitsecs=3600 